# resources/dist/init/spooter-server.service
# spooter-server.service

[Unit]
Description=Spooter Server
Documentation=https://github.com/thomasglachant/spooter
After=network.target

[Service]
Type=simple
User=spooter
Group=spooter
ExecStart=/usr/bin/spooter-server -config /etc/spooter/server.yaml
TimeoutStopSec=5s
LimitNOFILE=1048576
LimitNPROC=512
PrivateTmp=true
ProtectSystem=full
Restart=on-failure
RestartSec=5s

StandardOutput=append:/var/log/spooter/server.log
StandardError=append:/var/log/spooter/server.log

[Install]
WantedBy=multi-user.target
