# resources/dist/init/spooter-agent.service
# spooter-agent.service

[Unit]
Description=Spooter Agent
Documentation=https://github.com/thomasglachant/spooter
After=network.target

[Service]
Type=simple
User=spooter
Group=spooter
ExecStart=/usr/bin/spooter-agent -config /etc/spooter/agent.yaml
TimeoutStopSec=30s
LimitNOFILE=1048576
LimitNPROC=512
PrivateTmp=true
ProtectSystem=full
Restart=on-failure
RestartSec=5s

StandardOutput=append:/var/log/spooter-agent.log
StandardError=append:/var/log/spooter-agent.log

[Install]
WantedBy=multi-user.target
